# root/etc/s6-overlay/s6-rc.d/init-mod-kasm-rclone-plugin/run
ARCH=$(uname -m | sed 's/aarch64/arm64/g' | sed 's/x86_64/amd64/g')
DOCKER_PLUGIN_DIR=$(docker info | awk -F': ' '/Docker Root Dir/ {print $2}')-plugins

mkdir -p ${DOCKER_PLUGIN_DIR}/rclone/config
mkdir -p ${DOCKER_PLUGIN_DIR}/rclone/cache

docker plugin install rclone/docker-volume-rclone:${ARCH} args="-v" --alias rclone --grant-all-permissions
